# systemd/.config/systemd/user/mbsync.service
[Unit]
Description=Mailbox synchronization service

[Service]
Type=oneshot
ExecStart=/usr/bin/mbsync -V lucas
ExecStartPost=/usr/bin/notmuch new
