# systemd/user/restic-daily.timer
[Unit]
Description=Restic User Backup Timer
 
[Timer]
WakeSystem=false
OnCalendar=daily
Persistent=true
 
[Install]
WantedBy=timers.target
